# crontab/10.19.10.210.sh
#!/usr/bin/env bash

promtail_update_log

nvidia_metric
