# crontab/10.19.10.230.sh
#!/usr/bin/env bash

promtail_update_log


windowPost_metric

nvidia_metric
