# crontab/10.19.10.204.sh
#!/usr/bin/env bash

lotus_metric
miner_metric
