# crontab/10.9.7.173.sh
#!/usr/bin/env bash

promtail_update_log

lotus_metric

windowPost_metric

nvidia_metric
