# crontab/10.9.7.10.sh
#!/usr/bin/env bash

promtail_update_log
