# crontab/10.9.7.173.sh
#!/usr/bin/env bash

echo "${BASEURL}/promtail/update_log.sh"
# curl -SsL ${BASEURL}/promtail/update_log.sh | bash
