# crontab/10.9.7.176.sh
#!/usr/bin/env bash

lotus_metric

winningPost_metric

nvidia_metric
