# crontab/10.9.7.173.sh
#!/usr/bin/env bash

curl -SsL ${BASEURL}/promtail/update_log.sh | bash
