# crontab/10.9.7.175.sh
#!/usr/bin/env bash

lotus_metric

windowPost_metric

nvidia_metric
